# hostagent/Resources/LaunchAgents/com.haven.hostagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.haven.hostagent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/hostagent</string>
        <string>--config</string>
        <string>/Users/YOUR_USER/.haven/hostagent.yaml</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Library/Logs/Haven/hostagent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Library/Logs/Haven/hostagent-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HAVEN_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>Nice</key>
    <integer>0</integer>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
